# backup.crontab
# Ежедневное резервное копирование в 2:00
0 2 * * * /bin/bash /vagrant/backup.sh